# ansible-playbook/roles/docker-compose-systemd/templates/docker-compose-reload.service.j2
[Unit]
Description=Refresh images and update containers for {{ docker_compose_dir_path }}

[Service]
Type=oneshot

ExecStart=/bin/systemctl reload {{ docker_compose_service_unit_file_name }}