# media-server-ansible/roles/network-file-share-mount/templates/automount_unit.automount.j2
[Unit]
Description=Media Share Auto Mount

[Automount]
Where={{ local_mount_path }}

[Install]
WantedBy=multi-user.target